# tests/fixtures/tasks_config_template.conf
[program:sleep10]
command=sleep 10
numprocs=3
autorestart=true
process_name=%(program_name)s_%(process_num)02d

[program:sleep20]
command=sleep 20
numprocs=3
autorestart=true
process_name=%(program_name)s_%(process_num)02d

